# savara/customs/etc/local.d/crda.start
# /etc/conf.d/local.start

# temporary fix, how to do this correctly?
/usr/sbin/iw reg set CZ
