# savara/customs/etc/local.d/disable-wakeup.start

# disable USB and LAN wakup
echo XHC | sudo tee /proc/acpi/wakeup
echo GLAN | sudo tee /proc/acpi/wakeup



