# savara/customs/etc/local.d/microcode.start

# check if this is needed

echo 1 > /sys/devices/system/cpu/microcode/reload

