# zax/customs/var/cron/crontabs/root
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.XXXXKvs73e installed on Mon Oct 29 11:07:52 2018)
# (Cron version V5.0 -- $Id: crontab.c,v 1.12 2004/01/23 18:56:42 vixie Exp $)
MAILTO=""

# check scripts in cron.hourly, cron.daily, cron.weekly and cron.monthly
# Global variables
8,12,18,22  *  * * *  /etc/init.d/ntp-client restart


