# zax/customs/etc/local.d/gem_proxy.start

pkill cntlm

cntlm

export http_proxy=http://localhost:3129
export https_proxy=http://localhost:3129


