# zax/customs/etc/local.d/mountvbox.start

# anomen : mount VBox shared folders

mount -t vboxsf -o fmode=666 C_DRIVE /mnt/drive_c
mount -t vboxsf -o fmode=666 TC /mnt/tc

mount -t vboxsf -o fmode=666 DATASHARE /mnt/datashare
mount -t vboxsf -o fmode=666 DATASTORE /mnt/datastore

mount -t vboxsf -o fmode=666 TEMP /mnt/temp
