# zax/customs/etc/local.d/update_ip.start

# anomen : update windows hosts file


HOSTFILE=/mnt/drive_c/WINDOWS/system32/drivers/etc/hosts
IF=eth0

test -f "$HOSTFILE" || exit 1

cp -f "$HOSTFILE" /tmp/hosts.bak || exit 2
cp -f "$HOSTFILE" /tmp/hosts.tmp || exit 3

IP=`/sbin/ip addr show $IF | sed -n -e 's@.*inet\s*\([0-9.]*\)/.*@\1@p'`

sed -i -e "s/.*ZAX-ADDRESS-LINE.*/$IP zax zax.gemalto.com zax.cz www.zax.cz # ZAX-ADDRESS-LINE/"  /tmp/hosts.tmp

cp -f /tmp/hosts.tmp "$HOSTFILE"

